# systemd/metrics-collector.timer
[Unit]
Description=Run metrics collection every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=1s

[Install]
WantedBy=timers.target
